# conf/supervisor/jenkins.conf
[program:jenkins]
command = java -jar %(ENV_JENKINS_WAR)s  --webroot=/var/cache/jenkins/war --httpPort=%(ENV_JENKINS_PORT)s -Djava.awt.headless=true
autostart = true
startsecs = 10
user = jenkins
